# etc/profile.d/70-gke.sh
CLOUDSDK_CORE_PROJECT=dotnet-developer-experience
GKE_NAME=steeltoe-tap
GKE_REGION=us-east4
GKE_MACHINE=e2-standard-4
GKE_NODES=4
