# etc/profile.d/70-minukube.sh
MINIKUBE_CPUS=8
MINIKUBE_MEMORY=8192
MINIKUBE_DISK=8g
MINIKUBE_DRIVER=hyperkit
MINIKUBE_KUBERNETES_VERSION=1.21.1
